# sparse/etc/systemd/system/sensorfwd.service.d/override.conf
# For some reason sensorfwd doesn't exit in 15 seconds hanging the shutdown process,
# so kill it instead when stopping the service.
[Service]
ExecStop=/usr/bin/kill -9 $MAINPID
